# src/simplegraph/vindrik.service
[Unit]
Description=Vindriktning Simple Graph Logger

[Service]
Type=simple
ExecStart=/opt/vindrik/mqttrrd.py
Restart=on-failure
